# start.conf
start on startup
task
exec python /home/pi/Downloads/ech-pi-master/gpio-control.py